# services/wifipi.service
[Unit]
Description=Automatic startup for thermoSquid

[Service]
ExecStart=sudo hostapd /etc/hostapd/hostapd.conf

[Install]
WantedBy=default.target
